# modules/blue_team/ransomware_attack_chain/files/crontab
# Ransomware persistence
*/5 * * * * /tmp/encryptor --persist > /dev/null 2>&1
@reboot /tmp/encryptor --startup
# Hidden flag: FLAG{871ecac95cd88307}